# files/atboot.service
[Unit]
Description=BOINC node atboot tasks

[Service]
ExecStart=/usr/bin/sensors-detect --auto > /dev/null

[Install]
WantedBy=multi-user.target
