# files/atboot.service
[Unit]
Description=BOINC node atboot tasks

[Service]
ExecStart=/usr/local/bin/homefarm-atboot

[Install]
WantedBy=multi-user.target
